# dist/hyperpixel4-init.service
[Unit]
Description=HyperPixel 4.0" LCD Display Initialization
DefaultDependencies=no
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/hyperpixel4-init

[Install]
WantedBy=multi-user.target
